# elasticsearch5/templates/etc/cron.daily/remove_old_elasticsearch_logs.j2
#!/bin/sh

{{ ansible_managed | comment }}

find \
  /var/log/elasticsearch \
  -name '*.log.gz' \
  -mtime +{{ elasticsearch5_logs_keep_days }} \
  -exec rm --verbose {} \;
