# backups_rclone/templates/etc/systemd/system/backups_rclone@.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Script for regular backups based on rclone

[Service]
ExecStart=/opt/backups_rclone.sh %i
SyslogIdentifier=%p
Type=oneshot
