# rabbitmq/templates/etc/systemd/system/rabbitmq-server.service.d/limitnofile.conf.j2
{{ ansible_managed | comment }}

[Service]
LimitNOFILE={{ rabbitmq_limit_nofile }}
