# logstash/templates/etc/cron.daily/remove_old_logstash_logs.j2
#!/bin/sh

{{ ansible_managed | comment }}

find \
  /var/log/logstash \
  -name '*.log.gz' \
  -mtime +{{ logstash_logs_keep_days }} \
  -exec rm --verbose {} \;
