# haproxy/templates/etc/systemd/system/haproxy.service.d/syslogidentifier.conf.j2
{{ ansible_managed | comment }}

[Service]
SyslogIdentifier=haproxy-{{ _haproxy_major_version }}
